# agent/service-installers/launchd/com.cobweb.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.cobweb.agent</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/cobweb-agent</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>NetworkState</key>
        <true/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>5</integer>

    <key>StandardOutPath</key>
    <string>/var/log/cobweb-agent.log</string>
    <key>StandardErrorPath</key>
    <string>/var/log/cobweb-agent.log</string>

    <!-- root daemon — exec / file / pty all need it (see design.md §9). -->
    <key>UserName</key>
    <string>root</string>
    <key>GroupName</key>
    <string>wheel</string>

    <!-- Daemon defaults: aim at the EasyTier hub and trust the private CA.
         For local overrides drop a config.toml at /etc/cobweb-agent/. -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>RUST_LOG</key>
        <string>info</string>
        <key>COBWEB_AGENT_SERVER_URL</key>
        <string>wss://cobweb.lan:8088/agent/ws</string>
        <key>COBWEB_AGENT_TRUST_CA</key>
        <string>/etc/cobweb-agent/etmesh-ca.crt</string>
    </dict>

    <!-- macOS launchd has no native logrotate; if logs grow large, rotate
         via a separate periodic LaunchAgent. -->
</dict>
</plist>
